# src/daemon/deepin-anything-daemon.service
[Unit]
Description=Deepin anything daemon service
After=udisks2.service

[Service]
ExecStart=/usr/libexec/deepin-anything-daemon
Restart=on-failure
RestartSec=30
StartLimitIntervalSec=3600  # The time window is 3600 seconds (60 minutes)
StartLimitBurst=2           # A maximum of 2 starts are allowed within this window.
MemoryHigh=2G
MemoryMax=3G

[Install]
WantedBy=default.target
